# debian/postinst
# refresh 
systemctl daemon-reload
# start immediately
systemctl start radius-mac.service
# enable at boot
systemctl enable radius-mac.service

